# systemd/vix-arena.service.example
[Unit]
Description=VixArena C++ multiplayer sandbox
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=micu
Group=micu
WorkingDirectory=/home/micu/vix
EnvironmentFile=/home/micu/vix/.env
Environment=DATA_DIR=/var/lib/vix-arena
ExecStart=/home/micu/vix/build/vix-arena
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=read-only
ReadOnlyPaths=/home/micu/vix
StateDirectory=vix-arena
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
SystemCallArchitectures=native
SystemCallFilter=@system-service

[Install]
WantedBy=multi-user.target
