# systemd/vix-arena.service.example
[Unit]
Description=VixArena C++ multiplayer sandbox
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=micu
Group=micu
WorkingDirectory=/home/micu/vix
EnvironmentFile=/home/micu/vix/.env
ExecStart=/home/micu/vix/build/vix-arena
Restart=on-failure
RestartSec=3

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/micu/vix/data
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictNamespaces=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
RemoveIPC=true
UMask=0077
ProtectProc=invisible
ProcSubset=pid
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
